# cron.crontab
# m h dom mon dow command
@reboot /home/tim/listen/update.sh 3600
00 6 * * * /home/tim/listen/run.sh 3600