# cron.crontab
# m h dom mon dow command
00 6 * * * /home/tim/listen/upload.sh
15 6 * * * /home/tim/listen/run.sh 7200; sudo /usr/sbin/shutdown -h now
# 00 11 * * * /usr/sbin/shutdown -h now